# conf/redshift/.config/autostart/redshift.desktop

[Desktop Entry]
Exec=/usr/bin/redshift
Name=redshift
Comment=RedShift
Type=Application
