# systemd/.config/systemd/user/backup@.timer
[Unit]
Description=Clean Tarsnap Backups

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
