# systemd/.config/systemd/user/tarsnap-clean.service
[Unit]
Description=clean tarsnap backups
After=network-online.target

[Service]
ExecStart=/bin/sh -c tsar
KillSignal=SIGQUIT
