# systemd/.config/systemd/user/webcam-pic.timer
[Unit]
Description=Take a pic with the webcam every 5m

[Timer]
OnCalendar=*:0/10
Persistent=true

[Install]
WantedBy=timers.target

