# systemd/user/mail-zoho.path
[Unit]
Description=Sync mail on zoho maildir modification

[Path]
PathChanged=%h/mail/zoho/
Unit=mbsync@zoho.service


[Install]
WantedBy=default.target
