# systemd/user/webcam-pic.service
[Unit]
Description=Take a pic with the webcam

[Service]
ExecStart=/home/meatcar/bin/webcam-pic
