# conf/mbsync/.config/systemd/user/mail@.path
[Unit]
Description=Sync mail on %I maildir modification

[Path]
PathChanged=%h/mail/%i/
Unit=mbsync@%i.service

[Install]
WantedBy=default.target
